# host/aide/aide-check.timer
[Unit]
Description=AIDE weekly file integrity check timer
Documentation=man:aide

[Timer]
# Run weekly on Sunday at 02:30 UTC (= 04:30 SAST)
OnCalendar=Sun *-*-* 02:30:00 UTC
# Randomize by up to 15 minutes to avoid thundering herd
RandomizedDelaySec=900
# Persistent: catch up missed runs after boot
Persistent=true

[Install]
WantedBy=timers.target
