# host/aide/aide-check.service
[Unit]
Description=AIDE weekly file integrity check (databyte VPN portal)
Documentation=man:aide
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/aide-report.sh
User=root
Group=root
# Don't kill the process if it runs long — AIDE on large filesystems can take minutes
TimeoutStartSec=900
# Capture stdout/stderr in journald for later review
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aide-check
