# host/systemd/quota-monitor.service
[Unit]
Description=strongSwan quota monitor (Phase 6 — nft named meters + 80% warn + 100% hard cut)
Documentation=file:///home/zunaid/strongswan/quota/quota-monitor.py
After=docker.service quota-schema.service
Wants=quota-schema.service

# quota-monitor.py is the daemon that:
#   1. Samples per-VIP byte counters from nft NAMED METERS
#      (`ip filter client_src` for outbound, `client_dst` for inbound).
#      Meters are declared in /etc/nftables.conf (Phase 6, 2026-06-26);
#      kernel auto-creates counter elements per VIP on first packet.
#   2. Joins VIP → username via `docker exec strongswan swanctl --list-sas`
#   3. Joins username → customer via DB
#   4. Updates customers.data_used_bytes from counter delta
#   5. Fires 80% warn (logs + audit) on first crossing
#   6. Fires 100% hard cut: terminates active SAs, replaces EAP secret
#      in rw-eap.conf with KILLED-<random>, sets customers.over_quota=1
#
# Run on the VPS host (uses `docker exec strongswan` for VICI access).
# The script reads charon creds from /home/zunaid/strongswan/swanctl/conf.d/
# and writes back mutations there. StrongSwan container bind-mounts this dir.

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/zunaid/strongswan/quota/quota-monitor.py
Restart=on-failure
RestartSec=10
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
