# host/systemd/vpn-portal-smoke.service
[Unit]
Description=VPN portal + charon API-layer smoke test (Layer 4)
Documentation=https://github.com/Dippie-WP/databyte-Ikev2/blob/main/scripts/smoke.sh
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
# Run as root so docker exec + sqlite + /root file access works.
# Output → journald (read via journalctl -u vpn-portal-smoke)
ExecStart=/usr/local/bin/smoke.sh
# SuccessExitStatus: 0 = all OK, 1 = smoke failed (still "success" from systemd POV
# — the smoke script logs its own failures to journald + telegram).
SuccessExitStatus=0 1
# Don't auto-restart — this is a periodic check, not a service.