# host/cert-monitor/cert-expiry-check.service
[Unit]
Description=VPN portal TLS certificate expiry check
Documentation=man:openssl
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Output goes to journald; alerts read via journalctl -u cert-expiry-check
ExecStart=/usr/local/bin/cert-expiry-check.sh
# Exit code 2 = CRITICAL; surface via systemd + journal alerting
SuccessExitStatus=0 1
# Don't auto-restart (it's a check, not a service)
