# host/systemd/vpn-portal.service.d/readwrite-paths.conf
[Service]
# Allow portal to write its SSH known_hosts (StrictHostKeyChecking=accept-new needs to add host keys)
# /run/vpn-portal — socket dir
# /var/log/vpn-portal — gunicorn logs
# /var/lib/strongswan — SQLite DB (charon reads + portal reads via SSH)
# /var/lib/vpn-portal — portal SSH keys + known_hosts (NEW for VPS)
ReadWritePaths=/run/vpn-portal /var/log/vpn-portal /var/lib/strongswan /var/lib/vpn-portal
