# host/systemd/vpn-portal.service.d/runtime-dir.conf
[Service]
# Allow nginx (www-data) to traverse /run/vpn-portal/ to reach gunicorn.sock
# Socket file itself remains mode 0777; only directory traversal needs loosening
RuntimeDirectoryMode=0755
