# host/backup/backup-workspace.service
[Unit]
Description=Daily OpenClaw workspace backup to RustFS
Documentation=file:///root/projects/strongswan-vpn-gateway/host/backup/backup-workspace.sh
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Workspace dir is fixed; rclone config lives in /root/.config/rclone/rclone.conf
Environment=WORKSPACE_DIR=/root/.openclaw/workspace
# Run as root so rclone config is readable
ExecStart=/root/projects/strongswan-vpn-gateway/host/backup/backup-workspace.sh
# Capture stdout/stderr to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=backup-workspace
# 10min max — local copy should be <2s, 10min is a safety margin
TimeoutStartSec=600
