# host/backup/backup-vpn-portal-config.timer
[Unit]
Description=Daily VPN portal config + DB backup to RustFS
Documentation=file:///root/projects/strongswan-vpn-gateway/host/backup/backup-vpn-portal-config.service

[Timer]
# Daily at 03:30 UTC (05:30 SAST)
OnCalendar=*-*-* 03:30:00 UTC
Persistent=true
Unit=backup-vpn-portal-config.service

[Install]
WantedBy=timers.target
