# host/backup/backup-workspace.timer
[Unit]
Description=Daily OpenClaw workspace backup to RustFS
Documentation=file:///root/projects/strongswan-vpn-gateway/host/backup/backup-workspace.service

[Timer]
# Daily at 04:00 UTC (06:00 SAST). After VPN portal backup (03:30 UTC) to avoid simultaneous rclone load.
OnCalendar=*-*-* 04:00:00 UTC
Persistent=true
Unit=backup-workspace.service

[Install]
WantedBy=timers.target
